# deploy/cfc-formulaire.service
[Unit]
Description=CFC — Formulaire public d'auto-enregistrement des prospects
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/cfc-formulaire
# Variables (SYNC_TOKEN, PUBLIC_DB_PATH, PORT, RATE_SOUMISSION) :
EnvironmentFile=/opt/cfc-formulaire/.env
# 2 workers suffisent (volume faible, base SQLite mono-fichier).
ExecStart=/opt/cfc-formulaire/venv/bin/gunicorn -w 2 --threads 4 \
          -b 127.0.0.1:8080 app_public:app
Restart=always
RestartSec=3
# Durcissement
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/cfc-formulaire

[Install]
WantedBy=multi-user.target
